# scripts/sandbox-cleanup.service
[Unit]
Description=Whisper API Sandbox Cleanup
After=network.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/path/to/whisper-api
Environment=NODE_ENV=production
Environment=DATABASE_URL=your_database_url_here
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/node /path/to/whisper-api/scripts/sandbox-cleanup.js
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
